# ax5043/debian/lband.service
[Unit]
Description="AX5043 driver for L Band"

[Service]
Type=simple
ExecStart="/usr/bin/lband"
