# ax5043/debian/uhf.service
[Unit]
Description="AX5043 driver for UHF"

[Service]
Type=simple
ExecStart="/usr/bin/uhf"
